# .config/systemd/user/battery-notifier.timer
[Unit]
Description=Notify about low battery every minute

[Timer]
OnCalendar=*-*-* *:*:00
