# .config/systemd/user/battery-notifier.service
[Unit]
description=Notify about low battery

[Service]
ExecStart=%h/scripts/battery-notifier
