# .config/systemd/user/compton.service
[Unit]
Description=Compton compositor for X11
PartOf=graphical-session.target

[Service]
ExecStart=compton
