# .config/systemd/user/dwm-statusbar.service
[Unit]
Description=dwm status bar
PartOf=graphical-session.target

[Service]
ExecStart=%h/suckless/dwm/dwm_statusbar
Restart=on-failure
