# .config/systemd/user/dwm-statusbar.service
[Unit]
Description=dwm status bar
PartOf=graphical-session.target

[Service]
ExecStart=/home/phelipe/suckless/dwm/dwm_statusbar
