# .config/systemd/user/dwm-statusbar.service
[Unit]
Description=dwm status bar
PartOf=graphical-session.target

[Service]
ExecStart=%h/scripts/dwm-statusbar
Restart=on-failure
