# .config/systemd/user/dwm.service
[Unit]
Description=dwm
PartOf=graphical-session.target

[Service]
ExecStopPost=/bin/systemctl --user stop graphical-session.target
ExecStart=%h/.local/bin/dwm
Restart=on-failure
