# .config/systemd/user/nitrogen.service
[Unit]
Description=Restore previously set wallpaper
PartOf=graphical-session.target

[Service]
ExecStart=nitrogen --restore
