# .config/systemd/user/nm-applet.service
[Unit]
Description=Network Manager
PartOf=graphical-session.target

[Service]
ExecStart=nm-applet
