# .config/systemd/user/numlockx.service
[Unit]
Description=Activate Num Lock at startup
PartOf=graphical-session.target

[Service]
ExecStart=numlockx
