# .config/systemd/user/offlineimap.timer
[Unit]
Description=Sync offlineimap every hour
PartOf=graphical-session.target

[Timer]
OnCalendar=*-*-* *:00:00
