# .config/systemd/user/redshift.service
[Unit]
Description=redshift
PartOf=graphical-session.target

[Service]
ExecStart=redshift
