# .config/systemd/user/sxhkd.service
[Unit]
Description=Simple X Hotkey Daemon
PartOf=graphical-session.target

[Service]
ExecStart=sxhkd
ExecReload=kill -SIGUSR1 $MAINPID
