# templates/slock@.service.j2
[Unit]
Description=Lock X session using slock for user %i
Before=sleep.target

[Service]
User=%i
Environment=DISPLAY=:0
ExecStart=slock

[Install]
WantedBy=sleep.target
